# docs/systemd/xiaoban-agent.service.example
[Unit]
Description=Xiaoban-Agent My Stand runtime
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/xiaoban-agent
EnvironmentFile=-/etc/xiaoban-agent/xiaoban.env
Environment=XIAOBAN_HOME=/var/lib/xiaoban-agent
Environment=XIAOBAN_LOG_DIR=/var/log/xiaoban-agent
ExecStart=/opt/xiaoban-agent/.venv/bin/python -m xiaoban.cli gateway run --replace --accept-hooks
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
